# src/rc/RedHat/pdnsd.in
#!/bin/bash
#
# $Id: pdnsd.in,v 1.6 2001/07/26 12:26:55 tmm Exp $
#
#	/etc/rc.d/init.d/pdnsd
#
# Starts the Proxy DNS Daemon
#
# chkconfig: 2345 80 30
# description: Proxy DNS Daemon
# processname: pdnsd
# config: /etc/pdnsd.conf

prefix=@prefix@

test -x @prefix@/sbin/pdnsd || exit 0
[ -f @sysconfdir@/pdnsd.conf ] || exit 1

# Source function library.
. /etc/rc.d/init.d/functions

# Source networking configuration.
. /etc/sysconfig/network

# Check that networking is up.
[ ${NETWORKING} = "no" ] && exit 0

RETVAL=0

start() {
	echo -n 'Starting pdnsd: '
	daemon @prefix@/sbin/pdnsd -d
	RETVAL=$?
	echo
	[ $RETVAL -eq 0 ] && touch /var/lock/subsys/pdnsd
	return $RETVAL
}

stop() {
	echo -n 'Shutting down pdnsd: '
	killproc pdnsd
	RETVAL=$?
	echo
	[ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/pdnsd
	return $RETVAL
}

restart() {
	stop
	start
}

condrestart() {
	[ -e /var/lock/subsys/pdnsd ] && restart
	return 0
}

#
#	See how we were called.
#
case "$1" in
  start)
	start
	;;
  stop)
	stop
	;;
  status)
	status pdnsd
	;;
  reload|restart)
	restart
  	;;
  condrestart)
	condrestart
  	;;
  *)
	echo $"Usage: $0 {start|stop|status|restart|condrestart|reload}"
	RETVAL=1
esac

exit $RETVAL
